# high-noon.cron
0  12  * * * pi curl http://localhost:3000/playAudio/OW_high_noon.mp3
